# scripts/chimera.service
# =============================================================================
# Chimera HFT Systemd User Service
# =============================================================================
# Installation:
#   1. Copy to ~/.config/systemd/user/chimera.service
#   2. Enable: systemctl --user enable chimera
#   3. Start: systemctl --user start chimera
#   4. Check: systemctl --user status chimera
#   5. Logs: journalctl --user -u chimera -f
#
# For autostart on boot (even without login):
#   sudo loginctl enable-linger $USER
# =============================================================================

[Unit]
Description=Chimera HFT Trading Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Chimera/build
ExecStart=%h/Chimera/build/chimera
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/Chimera/logs/chimera_service.log
StandardError=append:%h/Chimera/logs/chimera_service.log

# Resource limits
LimitNOFILE=65535
LimitMEMLOCK=infinity

# CPU affinity (optional - uncomment to pin)
# CPUAffinity=1 2

[Install]
WantedBy=default.target
